# scripts/datajoint-ingest.cron
# run datajoint ingest (TODO: path hardcoding, flock)
45 */4 * * * /src/ibl-pipeline/datajoint-ingest.cron.sh
